# packet_snifferd.service
[Unit]
Description=Ip packet sniffer daemon
After=network.target

[Service]
ExecStart=packet_snifferd
