# packet_snifferd.service
[Unit]
Description=Ip packet sniffer daemon

[Service]
Type=simple
ExecStart=/usr/bin/packet_snifferd

