# packet_snifferd.service
[Unit]
Description=Ip packet sniffer daemon

[Service]
Type=dbus
BusName=snifferd.dest
ExecStart=/usr/bin/packet_snifferd


